# config/launch_agents/com.luabimlab.bimobject-crawl-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.luabimlab.bimobject-crawl-weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/opt/python@3.14/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
    <string>/Users/choejeong-yeon/LUA BIM LABS/scripts/crawl_bimobject.py</string>
    <string>--max-products</string>
    <string>150</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/choejeong-yeon/LUA BIM LABS</string>

  <!-- 매주 일요일 03:00 실행 — BIMobject 전체 카테고리 패밀리 목록 스캔 및 Obsidian 노트 업데이트 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/choejeong-yeon/LUA BIM LABS/logs/bimobject_crawl.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/choejeong-yeon/LUA BIM LABS/logs/bimobject_crawl.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/opt/python@3.14/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/choejeong-yeon</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
