# config/launch_agents/com.luabimlabs.sync-knowledge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.luabimlabs.sync-knowledge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/choejeong-yeon/LUA BIM LABS/scripts/sync_knowledge.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/choejeong-yeon/LUA BIM LABS/logs/sync_knowledge.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/choejeong-yeon/LUA BIM LABS/logs/sync_knowledge_err.log</string>
</dict>
</plist>
