# config/launch_agents/com.luabimlab.starter-telegram-sheet-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.luabimlab.starter-telegram-sheet-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/choejeong-yeon/LUA BIM LABS/.dev-venv/bin/python</string>
    <string>/Users/choejeong-yeon/LUA BIM LABS/scripts/starter_plan_sync_telegram_sheet.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/choejeong-yeon/LUA BIM LABS</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>600</integer>

  <key>StandardOutPath</key>
  <string>/Users/choejeong-yeon/LUA BIM LABS/runtime/starter_telegram_sheet_sync.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/choejeong-yeon/LUA BIM LABS/runtime/starter_telegram_sheet_sync.err</string>
</dict>
</plist>
