# scripts/com.luabimlabs.tailscale-userspace.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.luabimlabs.tailscale-userspace</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/opt/tailscale/bin/tailscaled</string>
    <string>--tun=userspace-networking</string>
    <string>--socket=/Users/choejeong-yeon/.local/run/tailscale/tailscaled.sock</string>
    <string>--statedir=/Users/choejeong-yeon/.local/share/tailscale</string>
    <string>--socks5-server=localhost:1055</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/choejeong-yeon/LUA BIM LABS/logs/tailscaled-userspace.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/choejeong-yeon/LUA BIM LABS/logs/tailscaled-userspace.err.log</string>
</dict>
</plist>
